# deployment/ctfrecon.service
[Unit]
Description=Gunicorn instance to serve CTF-Recon-Web
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/CTF-Recon-Web
Environment="PATH=/home/ubuntu/CTF-Recon-Web/venv/bin"
ExecStart=/home/ubuntu/CTF-Recon-Web/venv/bin/gunicorn --workers 3 --bind unix:ctfrecon.sock -m 007 app:app

[Install]
WantedBy=multi-user.target
